# richmenu-editor.service
[Unit]
Description=LINE Rich Menu Editor - Multi-user Collaboration
After=network.target

[Service]
Type=simple
User=your-user
Group=your-group
WorkingDirectory=/path/to/richmenu-editor/flask
Environment="PATH=/path/to/richmenu-editor/flask/venv/bin"
Environment="SECRET_KEY=change-this-to-random-secret-key"
Environment="PORT=1153"
Environment="DEBUG=False"
ExecStart=/bin/bash -c '/path/to/richmenu-editor/flask/venv/bin/gunicorn -k eventlet -w 1 -b 127.0.0.1:${PORT} app:app'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全性設定
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

